# deploy/com.constellation.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  cloudflared named tunnel `constellation-air` — system LaunchDaemon (Group 3).
  Runs the tunnel that fronts the loopback MCP origin (127.0.0.1:8000). The
  tunnel + credentials-file (<UUID>.json) are created in the morning runbook
  (cloudflared tunnel login/create); this daemon only *runs* an already-created
  tunnel. Config lives at deploy/cloudflared-config.yml (references the creds
  file by absolute path). Runs as `ama` so it reads ~/.cloudflared creds.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.constellation.tunnel</string>

    <key>UserName</key>
    <string>ama</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--no-autoupdate</string>
        <string>--config</string>
        <string>/Users/ama/dev/constellation-v3/deploy/cloudflared-config.yml</string>
        <string>run</string>
        <string>constellation-air</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TUNNEL_ORIGIN_CERT</key>
        <string>/Users/ama/.cloudflared/cert.pem</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ama/dev/constellation-v3/data/logs/tunnel.daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ama/dev/constellation-v3/data/logs/tunnel.daemon.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
